# update-dotfiles/.config/systemd/user/update-dotfiles.service
[Unit]
Description="update dotfiles"

[Service]
Type=oneshot
ExecStart=%h/bin/update-dotfiles.sh
