# update-dotfiles/.config/systemd/user/update-dotfiles.timer
[Unit]

[Timer]
OnCalendar=daily
Persistent=true
Unit=update-dotfiles.service

[Install]
WantedBy=timers.target
